# replica/etc/profile.d/shell-config.sh
#global shell config
alias ls='ls --color=auto'
alias ll='ls -lah --color=auto'
alias grep='grep --color=auto'
alias cls='clear'
alias poweroff='shutdown now'
alias reboot='reboot'
export YOUROS=1

